# tautulli/rootfs/etc/services.d/tautulli/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Tautulli
# Runs Tautulli
# ==============================================================================

# Start Tautulli
bashio::log.info "Starting Tautulli..."
exec \
  python /opt/Tautulli.py \
  --datadir /data \
  --nolaunch \
  -p 8181